# etc/systemd/netcat.service
[Unit]
Description=Netcat Listener Service
After=network.target

[Service]
ExecStart=/usr/bin/nc -l -p 8080
Restart=always
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
